# systemd/koios@120.service
[Unit]
Description=Koios Stake Distribution Update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/koios/stake-distribution-update.sh
