# systemd/koios@720.service
[Unit]
Description=Koios Pool Groups Update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/koios/pool-groups-update.sh
